# src/limit-service/db-sentry-limit.service
[Unit]
Description=dB-Sentry Limit Service
After=network.target mosquitto.service influxdb.service
Wants=network.target

[Service]
Type=simple
User=garth
Group=garth
WorkingDirectory=/home/garth/dB-Sentry/src/limit-service
Environment="PATH=/home/garth/dB-Sentry/src/limit-service/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/garth/dB-Sentry/src/limit-service/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give the service access to GPIO
SupplementaryGroups=gpio i2c spi

[Install]
WantedBy=multi-user.target
